# Assignment03-Code/hello-bash-service.service
[Unit]
Description=hello World Service

[Service]
Type=simple
ExecStart=/usr/local/bin/hello-bash-service.sh

[Install]
WantedBy=multi-user.target
